# distrib/build/sources/updates/update1/coretree.list
./httpd/html/help/dmzholes.cgi.html.en
./httpd/cgi-bin/advstatus.cgi
./httpd/cgi-bin/traffic.cgi
./httpd/cgi-bin/ids.cgi
./httpd/cgi-bin/logs.cgi/proxylog.dat
./httpd/cgi-bin/portfw.cgi
./httpd/cgi-bin/dmzholes.cgi
./httpd/cgi-bin/ddns.cgi
./httpd/cgi-bin/xtaccess.cgi
./usr/lib/smoothwall/p3scan.conf.in
./usr/lib/smoothwall/langs/base.pl
./usr/lib/smoothwall/smoothtype.pm
./usr/bin/smoothwall/writeipsec.pl
./usr/bin/smoothwall/writedhcp.pl
./usr/bin/smoothwall/writeproxy.pl
./usr/bin/smoothwall/machine_reg.pl
./etc/logrotate.d/smoothd
./etc/logrotate.conf
./etc/rc.d/rc.firewall.up
./etc/rc.d/rc.network
