# distrib/build/sources/coretree/tree/etc/cron.daily/bogons
#!/bin/bash

# Update the bogons list for IP Block
/usr/bin/perl /root/scripts/bogons.pl

