# distrib/build/sources/updates/update2/coretree.list
./httpd/cgi-bin/ids.cgi
./httpd/cgi-bin/logs.cgi/imviewer.cgi
./httpd/cgi-bin/portfw.cgi
./etc/crontab
./etc/inittab
./etc/cron.daily/runfreshclam
./etc/lilo.conf
./etc/rc.d/rc.network
