# distrib/build/sources/updates/update4/coretree.list
./httpd/html/help/im.cgi.html.en
./httpd/cgi-bin/dhcp.cgi
./httpd/cgi-bin/backup.img
./httpd/cgi-bin/logs.cgi/imviewer.cgi
./httpd/cgi-bin/im.cgi
./httpd/cgi-bin/ids.cgi
./httpd/cgi-bin/updates.cgi
./usr/lib/smoothwall/header.pm
./usr/lib/smoothwall/langs/base.pl
./usr/lib/smoothwall/smoothnet.pm
./usr/lib/smoothwall/imspector.conf.in
./usr/bin/smoothwall/writeim.pl
./usr/bin/smoothwall/writedhcp.pl
./usr/bin/smoothwall/writeproxy.pl
./usr/bin/smoothwall/imagerotator.pl
./usr/bin/smoothwall/setddns.pl
./usr/bin/smoothwall/writesiproxdconf.pl
./etc/rc.d/rc.sysinit
./etc/snort.conf
./usr/ssl/certs/
